# laravel-svelte-port/laravel/deploy/supervisor/nginx.ini
[program:nginx]
process_name=%(program_name)s
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stopwaitsecs=10
startsecs=5
startretries=3